# zsh/profile.d/00-path.zsh
add_to_path_if_not_exists "$HOME/.local/bin"
add_to_path_if_not_exists "$HOME/local/nvim/bin"
add_to_path_if_not_exists /opt/local/bin
add_to_path_if_not_exists /usr/local/bin
add_to_path_if_not_exists /usr/bin
add_to_path_if_not_exists /bin
add_to_path_if_not_exists /usr/sbin
add_to_path_if_not_exists /sbin
