# zsh/profile.d/10-mise.zsh
if command_exists mise; then
  eval "$(mise activate zsh)"
fi
